# sgk_test/templates.d/blackhole.service
[Unit]
Description=Blackhole service
After=network.target
[Service]
ExecStart="путь до скрипта"/blackhole
[Install]
WantedBy=multi-user.target
